# dock-tunnel/etc/profile.d/dock-vars

#######################################################################################################################
#
# A simple that initialized global variables for the dock.
#
# DOCK_GATEWAY     IP address or hostname of the dock
#
#######################################################################################################################

export DOCK_GATEWAY=$GATEWAY
